# .config/systemd/user/dark.timer
[Unit]
Description=Change the desktop appearance to a dark theme on sunset.

[Timer]
OnCalendar=*-*-* 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
